# magento/templates/supervisord-webserver.conf
[program:webserver]
command=/usr/bin/php -S 0.0.0.0:80 -t /data/pub/ /data/phpserver/router.php
autorestart=true